# ansible/roles/set_local_multicast/files/99-lo-multicast.rules
SUBSYSTEM=="net",ACTION=="add",KERNEL=="lo",RUN+="/usr/sbin/ip link set dev %k multicast on"